# etc/mbs-backup-to-cloud.timer
[Unit]
Description=Minecraft Bedrock Server Backup to Cloud Timer

[Timer]
OnCalendar=Fri *-*-* 6:00:00

[Install]
WantedBy=timers.target
