# etc/mbs-update.timer
[Unit]
Description=Minecraft Bedrock Server Update Timer

[Timer]
OnCalendar=Mon *-*-* 7:00

[Install]
WantedBy=timers.target
