# etc/mbs-update.timer
[Unit]
Description=Minecraft Bedrock Server Update Timer

[Timer]
OnCalendar=Sat 6:00:00

[Install]
WantedBy=timers.target
