# etc/mbs-update.service
[Unit]
Description=Minecraft Bedrock Server Update

[Service]
EnvironmentFile=/docker/minecraft-bedrock-server-setup/.env
User=${MY_USERNAME}
Group=${MY_USERNAME}
Type=oneshot
WorkingDirectory=${DIR_REPO}/bin
ExecStart=${DIR_REPO}/bin/update.sh

[Install]
WantedBy=timers.target
