# etc/mbs-update.service
[Unit]
Description=Minecraft Bedrock Server Update

[Service]
Type=oneshot
ExecStart=bash mbs/minecraft-bedrock-server-setup/bin/update.sh
