# etc/mbs-backup-cloud.service
[Unit]
Description=Minecraft Bedrock Server Backup to Cloud

[Service]
EnvironmentFile=/docker/minecraft-bedrock-server-setup/.env
User=${USERNAME}
Group=${USERNAME}
Type=oneshot
WorkingDirectory=${DIR_REPO}/bin
ExecStart=${DIR_REPO}/bin/rsync-googledrive.sh

[Install]
WantedBy=timers.target
