# etc/mbs-update.timer
[Unit]
Description=Minecraft Bedrock Server Update Timer

[Timer]
OnCalendar=Mon,Tue,Wed,Thu,Fri *-*-* 12:00:00

[Install]
WantedBy=timers.target
