# etc/mbs-backup.timer
[Unit]
Description=Minecraft Bedrock Server Backup Timer

[Timer]
OnCalendar=*-*-* 5:00

[Install]
WantedBy=timers.target
