# etc/mbs-backup.service
[Unit]
Description=Minecraft Bedrock Server Backup

[Service]
Type=oneshot
ExecStart=bash mbs/minecraft-bedrock-server-setup/bin/backup.sh
