# etc/mbs-backup-to-local.timer
[Unit]
Description=Minecraft Bedrock Server Backup to Local Timer

[Timer]
OnCalendar=Fri *-*-* 5:00:00

[Install]
WantedBy=timers.target
