# etc/mbs-backup-to-cloud.timer
[Unit]
Description=Minecraft Bedrock Server Backup to Cloud Timer

[Timer]
OnCalendar=Fri *-*-* 12:30:00

[Install]
WantedBy=timers.target
